# pool/bothook-provision.service
[Unit]
Description=BOTHook Provisioning Server (Baileys)
After=network-online.target
Wants=network-online.target

# If control-plane has taken over the WhatsApp login flow, do not run local provision server.
ConditionPathExists=!/opt/bothook/LOGIN_AUTHORITY.control-plane

[Service]
Type=simple
WorkingDirectory=/opt/bothook/provision
Environment=PROVISION_DATA_DIR=/opt/bothook/provision/data
Environment=PROVISION_PORT=18999
Environment=OPENCLAW_HOME=/home/ubuntu
Environment=OPENCLAW_STATE_DIR=/home/ubuntu/.openclaw
Environment=PATH=/home/ubuntu/.npm-global/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/usr/bin/node /opt/bothook/provision/server.mjs
Restart=always
RestartSec=2

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=false

[Install]
WantedBy=multi-user.target
