# p-site/artifacts/v0.2.23/systemd/openclaw-gateway.service
[Unit]
Description=OpenClaw Gateway (system-level)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
Environment=OPENCLAW_HOME=/home/ubuntu
Environment=OPENCLAW_CONFIG_PATH=/home/ubuntu/.openclaw/openclaw.json
Environment=OPENCLAW_GATEWAY_PORT=18789
# Ensure ubuntu user's npm-global bin is on PATH (used for chat-driven self-updates, too)
Environment=PATH=/home/ubuntu/.npm-global/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
WorkingDirectory=/home/ubuntu
ExecStart=/opt/bothook/bin/openclaw-gateway-start.sh
Restart=always
RestartSec=2

# Hardening (minimal; keep compatible)
# Allow sudo/root escalation when owner enables full-control mode.
# Note: sudo still requires /etc/sudoers.d rules; this only removes the Linux no_new_privs restriction.
NoNewPrivileges=false
PrivateTmp=true
ProtectSystem=full
ProtectHome=false

[Install]
WantedBy=multi-user.target
