# p-site/artifacts/v0.2.42/systemd/bothook-openclaw-watchdog.timer
[Unit]
Description=BOTHook Watchdog timer (openclaw high CPU mitigation)

[Timer]
OnBootSec=45s
OnUnitActiveSec=10s
AccuracySec=1s
Unit=bothook-openclaw-watchdog.service

[Install]
WantedBy=timers.target
