# p-site/artifacts/v0.2.7/systemd/bothook-postboot-verify.service
[Unit]
Description=BOTHook Post-boot Verification (P0.2 reboot acceptance)
After=network-online.target openclaw-gateway.service bothook-provision.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/bothook/bin/postboot_verify.sh
# Retry transient failures (e.g. gateway port not yet listening right after reboot)
Restart=on-failure
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=6

[Install]
WantedBy=multi-user.target
