# control-plane/systemd/bothook-welcome-gap-alert.service
[Unit]
Description=BOTHook alert: bound-but-no-welcome -> Telegram (deduped)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/.openclaw/workspace
Environment=BOTHOOK_DB_PATH=/home/ubuntu/.openclaw/workspace/control-plane/data/bothook.sqlite
Environment=BOTHOOK_OPENCLAW_BIN=/home/ubuntu/.npm-global/bin/openclaw
Environment=BOTHOOK_TELEGRAM_OWNER_CHAT_ID=7095719535
Environment=BOTHOOK_WELCOME_GAP_ALERT_DEDUPE_MIN=30
Environment=BOTHOOK_WELCOME_GAP_ALERT_LOOKBACK_MIN=20
ExecStart=/usr/bin/node /home/ubuntu/.openclaw/workspace/control-plane/workers/welcome_gap_alert.mjs
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
