# p-site/artifacts/v0.2.18/systemd/bothook-postboot-verify.service
[Unit]
Description=BOTHook Post-boot Verification (P0.2 reboot acceptance)
After=network-online.target openclaw-gateway.service bothook-provision.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/bothook/bin/postboot_verify.sh

# Prevent restart storms:
# - Postboot verify can be expensive (openclaw CLI + gateway restart).
# - If it fails repeatedly, it can pin CPU and destabilize WhatsApp onboarding.
# We still retry transient failures, but with much slower backoff and a lower burst.
Restart=on-failure
RestartSec=60
StartLimitIntervalSec=1800
StartLimitBurst=2

# Keep the system responsive even under load
Nice=10
IOSchedulingClass=idle

[Install]
WantedBy=multi-user.target
